# deploy/aapanel/supervisor-go-panel.conf
; aaPanel Supervisor config - Azure Panel Go runtime

[program:azure-panel-go]
command=/www/wwwroot/azure-panel/bin/azure-panel-go
directory=/www/wwwroot/azure-panel
user=www
autostart=true
autorestart=true
startsecs=2
startretries=3
stopwaitsecs=10
stdout_logfile=/www/wwwlogs/azure-panel-go.log
stderr_logfile=/www/wwwlogs/azure-panel-go-error.log
environment=GO_PANEL_HOST="127.0.0.1",GO_PANEL_PORT="3000",GO_PANEL_MODE="go",GO_PANEL_NODE_COMPAT_ENABLED="true",GO_PANEL_NODE_COMPAT_URL="http://127.0.0.1:3001",GO_PANEL_STATIC_DIR="/www/wwwroot/azure-panel/build/client",GO_PANEL_SUBMIT_DEADLINE_SECONDS="30",GO_PANEL_QUEUE_LIMIT="128"
