# deploy/aapanel/supervisor-worker.conf
; aaPanel Supervisor 配置 — 自动补机 Worker（独立进程）
; 与 Web 进程分离，避免 Nginx 重启或 Web 崩溃影响补机任务

[program:azure-panel-worker]
command=/usr/bin/node /www/wwwroot/azure-panel/build/worker.js
directory=/www/wwwroot/azure-panel
user=www
autostart=true
autorestart=true
startsecs=5
startretries=3
stopwaitsecs=10
stdout_logfile=/www/wwwlogs/azure-panel-worker.log
stderr_logfile=/www/wwwlogs/azure-panel-worker-error.log
environment=NODE_ENV="production"
