# deploy/aapanel/supervisor-go-replenisher.conf
; aaPanel Supervisor config - optional Go replenisher sidecar

[program:azure-panel-go-replenisher]
command=/www/wwwroot/azure-panel/bin/azure-panel-go-replenisher
directory=/www/wwwroot/azure-panel
user=www
autostart=true
autorestart=true
startsecs=2
startretries=3
stopwaitsecs=10
stdout_logfile=/www/wwwlogs/azure-panel-go-replenisher.log
stderr_logfile=/www/wwwlogs/azure-panel-go-replenisher-error.log
environment=GO_REPLENISHER_HOST="127.0.0.1",GO_REPLENISHER_PORT="43170",GO_REPLENISHER_MODE="observe",GO_REPLENISHER_SUBMIT_DEADLINE_SECONDS="30"
